# docker/mongodb/source/mongo-test.cron
0 5 31 2 * export AGE_KEY="/etc/mongodb-secret/age-key.txt" && /data/tools/mongo_manage.sh backup /etc/mongodb-secret/mongo.ini

